# system/overall/slackware/etc/cron.d/cron_fakehwclock.sh
#
# setting the time on boards without RTC
#
*/15 * * * * /usr/local/bin/fakehwclock.sh
